# health-tracker/health-tracker/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Redis..."
exec redis-server --port 6379 --bind 127.0.0.1 --dir /data/redis --save 300 10
