# health-tracker/health-tracker/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Health Tracker API..."
cd /opt/health-tracker/api
exec /opt/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
